# .docker/production/supervisord.conf
[supervisord]
nodaemon=true

[program:node]
command=node /home/app/server
